# srcs/images/influxdb/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:influxdb]
command=/usr/sbin/influxd &
autostart=true
autorestart=true

[program:telegraf]
command=/usr/bin/telegraf &
autostart=true
autorestart=true